# build/sysroot/usr/lib/systemd/system/vionetaos-local-storage-first.service
[Unit]
Before=docker.service
Description=Hello World

[Service]
ExecStart=/usr/bin/vionetaos-local-storage -init
Type=oneshot


[Install]
WantedBy=multi-user.target